# app/supervisor.d/php_fpm.conf
[program:php_fpm81]
command=php-fpm81 --nodaemonize
autorestart=true
priority=1

